# scripts/com.user.twstock-kline.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.user.twstock-kline</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/RONANFU/Downloads/test/twstock-kline/.venv/bin/python</string>
        <string>/Users/RONANFU/Downloads/test/twstock-kline/morning_brief.py</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/RONANFU/Downloads/test/twstock-kline</string>

    <!-- Run weekdays at 14:30 local time (TW market closes 13:30, give 1h buffer
         for Yahoo to publish the day's bar). -->
    <key>StartCalendarInterval</key>
    <array>
        <dict><key>Weekday</key><integer>1</integer><key>Hour</key><integer>14</integer><key>Minute</key><integer>30</integer></dict>
        <dict><key>Weekday</key><integer>2</integer><key>Hour</key><integer>14</integer><key>Minute</key><integer>30</integer></dict>
        <dict><key>Weekday</key><integer>3</integer><key>Hour</key><integer>14</integer><key>Minute</key><integer>30</integer></dict>
        <dict><key>Weekday</key><integer>4</integer><key>Hour</key><integer>14</integer><key>Minute</key><integer>30</integer></dict>
        <dict><key>Weekday</key><integer>5</integer><key>Hour</key><integer>14</integer><key>Minute</key><integer>30</integer></dict>
    </array>

    <key>StandardOutPath</key>
    <string>/Users/RONANFU/Downloads/test/twstock-kline/data/launchd.out.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/RONANFU/Downloads/test/twstock-kline/data/launchd.err.log</string>

    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
